# setup/robot.service
[Unit]
Description=Robot Runner

[Service]
ExecStart=/bin/python /home/$USER/quad/main.py

[Install]
WantedBy=multi-user.target